# rc.local
#!/bin/sh
# $FreeBSD$

/usr/sbin/pc-sysinstall -c /etc/pc-sysinstall.cfg
